# systemd/vortigaunt-monitor.service
[Unit]
Description=Vortigaunt System Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vortigaunt-monitor
