# systemd/vortigaunt-backup.service
[Unit]
Description=Vortigaunt Backup Service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vortigaunt-backup
