# systemd/vortigaunt-report.service
[Unit]
Description=Vortigaunt Health Report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vortigaunt-report
